# systemd/tgbot.service
# Example configuration for running the bot on a server as a systemd service
[Unit]
Description=YouTubeShortsDownloader daemon
After=network.target

[Service]
User=tgbot
Group=tgbot
Type=simple
WorkingDirectory=/opt/YouTubeShortsDownloader
ExecStart=/opt/YouTubeShortsDownloader/venv/bin/python /opt/YouTubeShortsDownloader/bot.py
Restart=always

[Install]
WantedBy=multi-user.target
